# setup-slaves/files/files/tmpfs.conf
[program:setup_tmpfs]
user=root
command=/bin/bash -c /root/tmpfs.sh
stdout_logfile=tmpfs_stdout.log
stderr_logfile=tmpfs_stderr.log
